# launch_scripts/stockpiler.timer
[Unit]
Description=Run Stockpiler Every 12 Hours

[Timer]
OnCalendar=*-*-* 00/12:00:00
Persistent=true

[Install]
WantedBy=timers.target
